# scheduler/crontab
SHELL=/bin/bash
BASH_ENV=/container.env

* * * * */2 root /usr/local/bin/python /app/periodic-tasks.py > /tmp/periodic.log
